# .ddev/web-build/drupal.cron
0 * * * * env $(grep -vE '^\s*#' /var/www/html/.ddev/.env.web | grep -vE '^\s*$' | xargs) /var/www/html/vendor/bin/drush cron
