# .ddev/web-build/drupal.cron
0 * * * * env $(grep -vE '^\s*#' "$DDEV_APPROOT/.ddev/.env.web" | grep -vE '^\s*$' | xargs) "$DDEV_APPROOT/vendor/bin/drush" cron --uri="$DDEV_PRIMARY_URL" --quiet
